# pkg/etc/init.d/wavefront-proxy
#!/bin/bash -e

# chkconfig: - 99 00
# description: Wavefront Proxy

### BEGIN INIT INFO
# Provides:          wavefront-proxy
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Required-Start:
# Required-Stop:
### END INIT INFO

#############################
service_name=$(basename "$0")
# Source custom settings
# shellcheck source=/dev/null
[[ -f "/etc/sysconfig/$service_name" ]] && . /etc/sysconfig/"$service_name"

desc=${DESC:-Wavefront Proxy}
user="wavefront"
wavefront_dir="/opt/wavefront"
proxy_dir=${PROXY_DIR:-$wavefront_dir/wavefront-proxy}
export JAVA_HOME="$proxy_dir/jre"
conf_file=${CONF_FILE:-$proxy_dir/conf/wavefront.conf}
log_file=${LOG_FILE:-/var/log/wavefront.log}
daemon_log_file=${DAEMON_LOG_FILE:-/var/log/wavefront-daemon.log}
err_file=$log_file
pid_file=${PID_FILE:-/var/run/$service_name.pid}
agent_jar=${AGENT_JAR:-$proxy_dir/bin/wavefront-push-agent.jar}
class="com.wavefront.agent.PushAgentDaemon"
app_args=${APP_ARGS:--f $conf_file}
java_args=${JAVA_ARGS:--XX:+AggressiveHeap}

# Legacy support for overrding JVM args in wavefront_proxy_launch.conf.
proxy_launch_conf="/opt/wavefront/wavefront-proxy/conf/wavefront_proxy_launch.conf"
if [[ -r $proxy_launch_conf ]]; then
	replacement_java_args=$(grep -ve '[[:space:]]*#' $proxy_launch_conf | tr '\n' ' ' | sed -e 's/^[[:space:]]*//' | sed -e 's/[[:space:]]*$//')
	if [[ "$replacement_java_args" != "-XX:+AggressiveHeap" ]]; then
		>&2 echo "!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!"
		>&2 echo "Using deprecated wavefront_proxy_launch.conf, this will go away in the future."
		>&2 echo "Please parameterize custom java arguments by overrding JAVA_ARGS in /etc/sysconfig/$service_name."
		>&2 echo "!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!!"
		java_args=$replacement_java_args
	fi
fi

jsvc=$proxy_dir/bin/jsvc

jsvc_exec()
{
	cd "$(dirname "$agent_jar")"
	[[ ! -f "$log_file" ]] && touch "$log_file" && chown "$user":"$user" "$log_file"
	# shellcheck disable=SC2086
	$jsvc \
		-user $user \
		-home $JAVA_HOME \
		-cp $agent_jar \
		$java_args \
		-outfile $log_file \
		-errfile $err_file \
		-pidfile $pid_file \
		-wait 20 \
		-debug \
		$1 \
		$class \
		$app_args &> $daemon_log_file
	if [[ $? -ne 0 ]]; then
		echo "There was a problem, see $err_file and $daemon_log_file"
	fi
}

case "$1" in
start)
	echo "Starting $desc"
	jsvc_exec
	echo "Done"
;;
status)
	if [[ -f "$pid_file" ]]; then
		echo "$desc is running (PID $(cat "$pid_file"))"
	else
		echo "$desc is not running."
		exit 3
	fi
;;
stop)
	echo "Stopping $desc"
	jsvc_exec "-stop"
	echo "Done"
;;
restart)
	if [[ -f "$pid_file" ]]; then
		echo "Stopping $desc"
		jsvc_exec "-stop"
	fi
	echo "Starting $desc"
	jsvc_exec
	echo "Done"
;;
*)
	echo "Usage: $0 {status | start | stop | restart}"
	exit 1
esac
